# includes.container/etc/xdg/autostart/org.vanillaos.fixme.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/fixme
Name=Fixme
